# to-machine/of-vagrant/by-ansible/roles/worker/files/concourse-worker.service
[Unit]
Description=Concourse CI worker process

[Service]
User=root
Restart=on-failure
EnvironmentFile=/etc/concourse/worker_environment
ExecStart=/usr/local/bin/concourse worker

[Install]
WantedBy=multi-user.target
